# loupedeck.service.template
[Unit]
Description=Loupedeck Linux Controller (User Session)
Documentation=https://github.com/archfill/loupedeck-linux
After=graphical-session.target sound.target

[Service]
Type=simple
WorkingDirectory={{PROJECT_DIR}}
Environment=NODE_ENV=production
Environment=LOG_LEVEL=info
ExecStart={{PROJECT_DIR}}/scripts/start-loupedeck.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
